# ops/scheduler/systemd/toastit-backup.service
[Unit]
Description=Toastit MariaDB backup to GCS
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/etc/toastit/scheduler.env
WorkingDirectory=/home/debian/toastit
ExecStart=/home/debian/toastit/ops/scheduler/backup-to-gcs.sh
User=debian
Group=debian
